# rootfs/etc/cron.d/tabby-cron
# This helps run the cron.php script for Tabby

0 5 * * * php '/config/tabby/cron.php'